# debian/postinst
#!/bin/bash

set -e

if [ "$1" = "configure" ]; then
    echo "Configuring biubo-waf..."
    
    if ! getent passwd biubo > /dev/null 2>&1; then
        echo "Creating biubo user..."
        useradd --system --home /etc/biubo-waf --no-create-home --shell /usr/sbin/nologin --group biubo
    fi
    
    mkdir -p /etc/biubo-waf
    mkdir -p /var/lib/biubo-waf
    mkdir -p /var/log/biubo-waf
    
    chown -R biubo:biubo /etc/biubo-waf
    chown -R biubo:biubo /var/lib/biubo-waf
    chown -R biubo:biubo /var/log/biubo-waf
    
    chmod 750 /etc/biubo-waf
    chmod 750 /var/lib/biubo-waf
    chmod 750 /var/log/biubo-waf
    
    if ! systemctl is-enabled biubo-waf >/dev/null 2>&1; then
        systemctl enable biubo-waf >/dev/null 2>&1 || true
    fi
    
    echo "Configuration complete."
fi
